# rootfs/etc/services.d/worker-beat/run
#!/command/with-contenv bash
set -euo pipefail

# shellcheck source=/dev/null
. /opt/dify-aio/lib/env.sh
load_generated_env
configure_dify_env

wait_for_http_url "http://127.0.0.1:5001/health" "Dify API"

export MODE=beat
export MIGRATION_ENABLED=false

cd /app/api
exec /entrypoint.sh
